# pulseaudio-autostart-fix/etc/xdg/autostart/pulseaudio-autostart-fix.desktop
[Desktop Entry]
Type=Application
Name=pulseaudio-autostart-fix
StartupNotify=false
Exec=systemctl --user restart pulseaudio.socket
